# ansible/templates/openldap-backup.service.j2
[Unit]
Description=OpenLDAP slapcat backup
Documentation=man:slapcat(8)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openldap-backup
